# systemd/twitch-fax.service
[Unit]
Description=Twitch FAX Service
After=network.target bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
User=%i
# bluetoothグループが存在する場合のみ使用
SupplementaryGroups=bluetooth
WorkingDirectory=/home/%i/twitch-fax/dist
ExecStart=/home/%i/twitch-fax/dist/twitch-fax
Restart=on-failure
RestartSec=10

# Bluetooth権限を付与
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/twitch-fax/dist

# 環境変数ファイル（オプション）
EnvironmentFile=-/home/%i/twitch-fax/dist/.env

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twitch-fax

[Install]
WantedBy=multi-user.target