# systemd/twitch-overlay.service
[Unit]
Description=Twitch Overlay Service
After=network.target bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
User=%i
# bluetoothグループが存在する場合のみ使用
SupplementaryGroups=bluetooth
WorkingDirectory=/home/%i/twitch-overlay/dist
ExecStart=/home/%i/twitch-overlay/dist/twitch-overlay
# 再起動ポリシー: エラー時と終了コード75で再起動
Restart=on-failure
RestartExitStatus=75
RestartSec=10

# Bluetooth権限を付与（強化）
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
# システムファイルへの書き込みを制限しつつ、必要なディレクトリへのアクセスは許可
ProtectSystem=false
ProtectHome=false

# Bluetoothデバイスへのアクセスを許可
PrivateDevices=no
DevicePolicy=auto

# 環境変数設定
Environment="SERVER_PORT=3456"
Environment="RUNNING_AS_SERVICE=true"
Environment="TWITCH_OVERLAY_DATA_DIR=%h/.twitch-overlay"

# 環境変数ファイル（オプション）
EnvironmentFile=-/home/%i/twitch-overlay/dist/.env

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twitch-overlay

[Install]
WantedBy=multi-user.target